# app-scientific/geant4/autobuild/overrides/etc/profile.d/geant4.sh
# Resolve OpenGL compatibility issues with some GPUs by using a specified graphics system
# Ref: https://geant4-userdoc.web.cern.ch/UsersGuides/AllGuides/html/ForApplicationDevelopers/Visualization/visdrivers.html
export G4VIS_DEFAULT_DRIVER=TSG
