# app-network/kea/autobuild/overrides/usr/lib/systemd/system/kea-dhcp4-server.service
[Unit]
Description=Kea IPv4 DHCP server
After=network.target

[Service]
Type=exec
User=kea
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
ExecStart=/usr/bin/kea-dhcp4 -c /etc/kea/kea-dhcp4.conf
KillSignal=SIGINT
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
